# deploy/claimbot.service
[Unit]
Description=ClaimBot Discord Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/ClaimBot
EnvironmentFile=/home/ubuntu/ClaimBot/.env
ExecStart=/usr/bin/python3 bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
